# services/stunnel.service
[Unit]
Description=SSL tunnel for network daemons
After=network.target

[Service]
ExecStart=/usr/bin/stunnel /etc/stunnel/stunnel.conf
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
Type=forking
PIDFile=/var/run/stunnel4.pid

[Install]
WantedBy=multi-user.target
